# code/setup/config/rc.local
#!/bin/bash
# To be copied to /etc/rc.local

modprobe dwc2
modprobe gadgetfs

mkdir /dev/gadget
mount -t gadgetfs gadgetfs /dev/gadget

# TODO: make a uMTP while-true loop in a separate script
/usr/bin/umtprd &
python3 /home/pi/client/__init__.py
exit 0
